# com.ledgeai.scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ledgeai.scraper</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/masak/Desktop/ニューススクレイピング/run_scraper.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/masak/Desktop/ニューススクレイピング</string>
    
    <key>StartInterval</key>
    <integer>10800</integer>
    
    <key>StandardOutPath</key>
    <string>/Users/masak/Desktop/ニューススクレイピング/logs/scraper.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/masak/Desktop/ニューススクレイピング/logs/scraper_error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
