# com.ainews.webapp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ainews.webapp</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/masak/Desktop/ニューススクレイピング/run_webapp.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/masak/Desktop/ニューススクレイピング/webapp</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/masak/Desktop/ニューススクレイピング/logs/webapp.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/masak/Desktop/ニューススクレイピング/logs/webapp_error.log</string>
</dict>
</plist>
